# bin/crontab
0      */1      *       *       *       /opt/infoservant.com/bin/backup_db.sh
