# bin/feeder.002.conf
# Start when system enters runlevel 2 (multi-user mode).
start on runlevel 3

env DBI_DSN="dbi:Pg:dbname=scotch_egg"

# Start delayed_job via the daemon control script.
exec /opt/infoservant.com/bin/feeder.pl

# Restart the process if it dies with a signal
# or exit code not given by the 'normal exit' stanza.
respawn

# Give up if restart occurs 10 times in 90 seconds.
respawn limit 10 90
